# test/featuretests/images/zookeeper/supervisor_zookeeper.conf
[program:zookeeper]
command=bash /var/lib/dist/zookeeper/start.sh
environment=ZOO_LOG_DIR=/var/log/zookeeper,ZOO_LOG4J_PROP='INFO,ROLLINGFILE'
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
priority=5
